# meta-core/recipes-graphics/telechips-egl/tc-teapot-light/tc-teapot-light.path
[Unit]
Description=Telechips EGL Sample Application path when weston ready

[Path]
PathExists=/var/run/user/root/wayland-0

[Install]
WantedBy=multi-user.target
