# meta-extra/recipes-benchmark/glmark2/glmark2/glmark2.path
[Unit]
Description=This service do print thermal temperature for current measuring

[Path]
PathExists=/var/run/user/root/wayland-0

[Install]
WantedBy=multi-user.target
